# patra-infra/scripts/dev.patra.nacos-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<!--
  Patra Nacos Tunnel - launchd agent for MacBook.

  Why this exists:
  - Patra 基础设施容器（含 Nacos 3.x）部署在 mac mini 上，MacBook 跨 tailscale 访问
  - tailscale wireguard 隧道 MTU 1280，Nacos gRPC server 启动后的 HTTP/2 SETTINGS
    frame 经 Docker bridge (MTU 1500) → OrbStack vpnkit → macOS host → tailscale
    时，frame 大小超过 1280 字节被静默丢弃，nacos-client 永远收不到 SETTINGS ACK
    导致 gRPC handshake 卡死（关联 tailscale issue #311 PMTU discovery 未实现）
  - 通过 ssh -L tunnel 绕过：ssh channel 在应用层重新分段 TCP 字节流，每段 ≤ 1280

  What it does:
  - 后台维持 SSH 连接到 mac mini，转发 8848 (HTTP) / 9848 (gRPC) / 8080 (console)
    到 MacBook 本地 127.0.0.1 同端口
  - keepalive 30s + ExitOnForwardFailure，连接断开时 launchd 自动重启进程

  应用配置（spring application yml）：
    spring.cloud.nacos.discovery.server-addr: ${NACOS_HOST:127.0.0.1}:${NACOS_PORT:8848}
    （NACOS_HOST 默认走 tunnel，无需任何额外环境变量）

  Install / uninstall: bash patra-infra/scripts/install-nacos-tunnel.sh {install|uninstall|status}
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.patra.nacos-tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-N</string>
        <string>-T</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>BatchMode=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=accept-new</string>
        <string>-L</string>
        <string>8848:127.0.0.1:8848</string>
        <string>-L</string>
        <string>9848:127.0.0.1:9848</string>
        <string>-L</string>
        <string>8080:127.0.0.1:8080</string>
        <string>linqibin@linqibins-mac-mini</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/patra-nacos-tunnel.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/patra-nacos-tunnel.err.log</string>
</dict>
</plist>
